# packaging/arch/spdif-keepalive.install
post_install() {
  cat <<'EOF'
spdif-keepalive installed.

Enable it for your user with:
  systemctl --user enable --now spdif-keepalive.service

Optional config:
  ~/.config/spdif-keepalive/config.toml

List sinks:
  spdif-keepalive --list-sinks
EOF
}

post_upgrade() {
  post_install
}
